# systemd/disk-cleaner.service
[Unit]
Description=Disk Cleaner — scheduled cache cleanup

[Service]
Type=oneshot
ExecStart=%h/.local/bin/disk-cleaner scheduled
